# deploy/worldcap-wa.timer
[Unit]
Description=worldcap WhatsApp daily post timer

[Timer]
OnCalendar=*-*-* 09:05:00
Persistent=true
Unit=worldcap-wa.service

[Install]
WantedBy=timers.target
